# config/supervisord.conf
# =============================================================================
# VibeOS Supervisor Configuration
# =============================================================================
#
# This file defines all services that run inside the VibeOS container.
# Services are started in priority order (lower = earlier).
#
# Service Stack:
#   1. xvfb      (100) - Virtual framebuffer (X11 display :0)
#   2. openbox   (150) - Minimal window manager
#   3. x11vnc    (200) - VNC server exposing the display
#   4. novnc     (300) - WebSocket proxy + web client (port 6080)
#   5. opencode  (50)  - AI coding assistant API (port 4096)
#
# Access Points:
#   - VNC:    localhost:5900 (raw VNC protocol)
#   - noVNC:  http://localhost:6080/vnc.html (default UI)
#             http://localhost:6080/beta.html (VibeOS custom UI)
#   - API:    http://localhost:4096 (OpenCode REST API)
#
# Logs: /var/log/supervisor/*.log
# =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# =============================================================================
# OpenCode AI Server
# =============================================================================
# Headless AI coding assistant that provides:
#   - REST API for session management
#   - SSE events for streaming responses
#   - Tool execution (bash, file ops, etc.)
# The Electron shell-ui connects to this via localhost:4096
[program:opencode]
command=/usr/bin/opencode serve --port 4096 --hostname 0.0.0.0
user=vibe
directory=/home/vibe
autostart=true
autorestart=true
priority=50
startsecs=3
startretries=5
stdout_logfile=/var/log/supervisor/opencode.log
stderr_logfile=/var/log/supervisor/opencode.err
environment=HOME="/home/vibe",USER="vibe",DISPLAY=":0"

# =============================================================================
# Display Stack
# =============================================================================

# Xvfb - Virtual Framebuffer
# Creates a virtual X11 display without physical hardware
# Resolution is set via RESOLUTION env var (default: 1920x1080)
[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 %(ENV_RESOLUTION)sx24
user=root
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err

# Openbox - Minimal Window Manager
# Lightweight WM that provides window decorations and management
# Config: /home/vibe/.config/openbox/
# Autostart: /home/vibe/.config/openbox/autostart (launches shell-ui)
[program:openbox]
command=/usr/bin/openbox-session
user=vibe
autostart=true
autorestart=true
priority=150
startsecs=2
stdout_logfile=/var/log/supervisor/openbox.log
stderr_logfile=/var/log/supervisor/openbox.err
environment=HOME="/home/vibe",USER="vibe",DISPLAY=":0"

# =============================================================================
# Remote Access Stack
# =============================================================================

# x11vnc - VNC Server
# Exposes the virtual display via VNC protocol on port 5900
# Flags:
#   -nopw      No password (container internal use)
#   -forever   Don't exit when client disconnects
#   -shared    Allow multiple simultaneous connections
[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -nopw -forever -shared -rfbport 5900 -xkb -noxrecord -noxfixes -noxdamage -noshm
user=root
autostart=true
autorestart=true
priority=200
startsecs=3
startretries=5
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.err
environment=DISPLAY=":0"

# noVNC - WebSocket Proxy + Web Client
# websockify bridges WebSocket (port 6080) to VNC (port 5900)
# Also serves static files from /opt/novnc/ including:
#   - vnc.html     Default noVNC interface
#   - beta.html    Custom VibeOS interface
#   - vnc_lite.html  Minimal reference implementation
[program:novnc]
command=websockify --web=/opt/novnc/ 6080 localhost:5900
user=vibe
autostart=true
autorestart=true
priority=300
startsecs=5
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err
environment=HOME="/home/vibe",USER="vibe"

# =============================================================================
# Service Group
# =============================================================================
# All services are grouped together for easy management:
#   supervisorctl start vibeos:*
#   supervisorctl stop vibeos:*
#   supervisorctl status vibeos:*
[group:vibeos]
programs=opencode,xvfb,openbox,x11vnc,novnc
priority=999
